# TestSkripte_Unterlagen/effekta_pv_anlage.service
[Unit]
Description=effekta_pv_anlage

[Service]
Type=simple
ExecStart=python3 /home/pi/effektaBmsMqtt.py

[Install]
WantedBy=multi-user.target
